# device-config/rootfs/etc/systemd/system/mypsplash-start.service
[Unit]
Description=Psplash Boot screen
After=systemd-vconsole-setup.service

[Service]
ExecStart=/usr/bin/psplash

[Install]
WantedBy=sysinit.target
